# handheld/rfid-autosync.service
[Unit]
Description=RFID Database Auto-Sync Service
After=network-online.target
Wants=network-online.target
Restart=always

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Integrated-RFID-Reader/handheld
ExecStart=/usr/bin/python3 /home/pi/Integrated-RFID-Reader/handheld/auto_sync_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/Integrated-RFID-Reader/handheld

[Install]
WantedBy=multi-user.target